# etc/get-cms-user-requests.cron
*/10 * * * * root [[ ! -f /usr/libexec/htcondor-pilot-job-router/get-cms-user-requests ]] || /usr/libexec/htcondor-pilot-job-router/get-cms-user-requests
