# etc/gums-grid-mapfile-local.cron
0 */6 * * * root [[ ! -f /usr/libexec/htcondor-pilot-job-router/gums-grid-mapfile-local ]] || /usr/libexec/htcondor-pilot-job-router/gums-grid-mapfile-local
